# systemd/routergeist.service
[Unit]
Description=Router Geist Backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=routergeist
Group=routergeist
EnvironmentFile=/opt/routergeist/.env
WorkingDirectory=/opt/routergeist/backend
ExecStart=/opt/routergeist/backend/start.sh
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
AmbientCapabilities=
CapabilityBoundingSet=
LockPersonality=true
RestrictNamespaces=true
RestrictRealtime=true
SystemCallFilter=@system-service

[Install]
WantedBy=multi-user.target


